# supervisord/supervisord.conf
[supervisord]

[supervisorctl]

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /app/uwsgi/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
